# systemd/cam-supervisor.service
[Unit]
Description=camrig capture supervisor (camera owner, scheduler, Cloudflare link)
After=network.target time-sync.target
Wants=network.target

[Service]
Type=exec
# Keep the user's own primary group: the secrets in /etc/camrig are
# root:$CAM_USER 0640 and unreadable if Group= overrides it with 'video'.
User=spaia
Group=spaia
SupplementaryGroups=video render
ExecStart=/opt/camrig/venv/bin/camrig supervise
Restart=on-failure
RestartSec=5
# Capture is the rig's whole job; let it use the CPU it needs.
Nice=-5

[Install]
WantedBy=multi-user.target
